# debian/dropbear.postinst
#!/bin/sh
set -e

test "$1" = 'configure' || exit 0

if test -x /etc/init.d/dropbear; then
  update-rc.d dropbear start 20 2 . stop 20 0 6 . >/dev/null
fi
